# ansible/roles/nexus/templates/etc/profile.d/nexus.sh
NEXUS_HOME="/usr/lib/nexus/nexus-{{nexus_version}}"
