# ansible/roles/vault/templates/etc/profile.d/vault.sh
export VAULT_ADDR=http://vault.service.consul:{{vault_port}}
